# foreman.target
[Unit]
Wants= foreman-web.target foreman-api.target

[Install]
WantedBy=multi-user.target
