# foreman-web.target
[Unit]
Requires=foreman.target
Wants= foreman-web-1.service
